# systemd/bg-removal.service
[Unit]
Description=Background Removal API (IS-Net + SAM ViT-L)
After=network.target
# Start after vLLM to ensure GPU is initialized
After=vllm.service

[Service]
Type=simple
User=amello
WorkingDirectory=/home/amello/bg-removal-service
Environment="PATH=/home/amello/bg-removal-env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/amello/bg-removal-env/bin/uvicorn api:app \
    --host 0.0.0.0 \
    --port 8002 \
    --workers 1 \
    --log-level info
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
